# deploy_tools/gunicorn-upstart-template.conf
description "Gunicorn server for SITENAME"

start on net-device-up
stop on shutdown

respawn

setuid USER
chdir SOURCEDIR

exec RELATIVE/PATH/TO/VIRTUALENV --bind unix:/tmp/SITENAME.socket PROJECT.wsgi:application
